# src/app/ext/servers/crontab
# use imega/supercronic container
* * * * * hx --unix-socket /var/www/sockets/local-web.socket GET http://web/minute/ext > /dev/null
